# airootfs/etc/skel/.config/autostart/wireplumber.desktop
[Desktop Entry]
Name=Wireplumber
Comment=Pipewire Session Manager
Exec=wireplumber
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
