# src/monody-desktop-config/etc/skel/.config/autostart/pipewire.desktop
[Desktop Entry]
Name=Pipewire
Comment=Pipewire Multimedia Server
Exec=pipewire
Terminal=false
Type=Application
X-GNOME-Autostart-enabled=true
