# src/monody-desktop-config/etc/skel/.config/autostart/xcape.desktop
[Desktop Entry]
Type=Application
Name=xcape
Exec=xcape -e 'Super_L=Alt_L|F1'
StartupNotify=false
Terminal=false
Hidden=false
