# airootfs/etc/runit/sv/agetty-autologin-tty1/run
#!/bin/sh
exec agetty --autologin live --noclear tty1 linux
